# rorolite/files/etc/profile.d/rorolite.sh
export PATH=/opt/rorolite/project/.rorolite/env/bin:$PATH