# rorolite/files/etc/profile.d/rorolite.sh
export PATH=/opt/rorolite/project/.rorolite/env/bin:/opt/anaconda3/bin:$PATH